# .config/udev_rules/52-bluetooth.rules
ACTION=="add", SUBSYSTEM=="input", ATTR{id/vendor}=="0067", ATTR{id/product}=="24a7", RUN+="/home/anson/.config/bluetooth_hfp.sh"
